# ops/systemd/veltrix-sepolia.service
[Unit]
Description=Veltrix Sepolia OP Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/home/ec2-user/Veltrix
ExecStart=/home/ec2-user/Veltrix/ops/start-sepolia-stack.sh
ExecStop=/usr/bin/docker compose -p veltrix-sepolia -f /home/ec2-user/Veltrix/ops/docker-compose.sepolia.yml down
RemainAfterExit=yes
TimeoutStartSec=0
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
